# config/kiwix.conf
[program:kiwix]
command=/var/kiwix-serve --port 8181 --library /library.xml
autostart=true
autorestart=true
stderr_logfile=/var/log/kiwix.err.log
stdout_logfile=/var/log/kiwix.out.log
