# config/kolibri.conf
[program:kiwix]
command=/usr/local/bin/kolibri start
autostart=true
autorestart=true
stderr_logfile=/var/log/kolibri.err.log
stdout_logfile=/var/log/kolibri.out.log
